# mimo-code/rootfs/etc/s6-overlay/s6-rc.d/mimocode/run
#!/usr/bin/with-contenv bashio
# ============================================================
# MiMo Code AI Server — S6 Service Entry
# Managed by s6-overlay, auto-restarts on crash
# ============================================================

PORT=$(bashio::config 'port')

exec /usr/local/bin/mimo serve --port "${PORT}"