# mimo-code/rootfs/etc/s6-overlay/s6-rc.d/mimocode-webui/run
#!/usr/bin/with-contenv bashio

CONFIG_PATH=/data/options.json
if [ -f "$CONFIG_PATH" ]; then
    MIMO_PORT=$(python3 -c "import json; print(json.load(open('$CONFIG_PATH')).get('port', 14096))")
else
    MIMO_PORT=14096
fi

# The web UI talks to the inner (127.0.0.1) port that mimo serve listens on
INNER_PORT=$((MIMO_PORT - 1))

# Persist channel config to /data (survives add-on updates), not the image layer
mkdir -p /data/mimocode

# Bridge add-on options (Feishu/WeChat) into the env vars server.py expects.
# server.py's _build_config_from_env() only reads FEISHU_*/WECHAT_*; without
# this bridge the credentials in /data/options.json never reach the server and
# the channels never connect.
# Note: using python3 instead of bashio::config because the container may not
# have a functional with-contenv / s6 context when started outside Supervisor.
export FEISHU_ENABLED=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(str(d.get('feishu_enabled',False)).lower())")
export FEISHU_APP_ID=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('feishu_app_id',''))")
export FEISHU_APP_SECRET=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('feishu_app_secret',''))")
export WECHAT_ENABLED=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(str(d.get('wechat_enabled',False)).lower())")
export WECHAT_CORP_ID=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('wechat_corp_id',''))")
export WECHAT_AGENT_ID=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('wechat_agent_id',''))")
export WECHAT_SECRET=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('wechat_secret',''))")
export WECHAT_TOKEN=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('wechat_token',''))")
export WECHAT_ENCODING_AES_KEY=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('wechat_encoding_aes_key',''))")
export PERSONAL_WECHAT_ENABLED=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(str(d.get('personal_wechat_enabled',False)).lower())")
export HA_MCP_URL=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('ha_mcp_url',''))")

# Overlay deployed code from the persistent volume /data/mimocode/webui.
# deploy_reliable.sh writes the full webui there so code survives `ha addons update`
# (which rebuilds the container from the image and would otherwise wipe docker-cp'd files).
# To roll back a bad deploy: `docker exec $CONTAINER rm -rf /data/mimocode/webui`.
if [ -d /data/mimocode/webui ] && [ -n "$(ls -A /data/mimocode/webui 2>/dev/null)" ]; then
    cp -rf /data/mimocode/webui/. /usr/share/mimocode/webui/
fi

exec env \
    MIMO_PORT="${INNER_PORT}" \
    MIMO_CONFIG="/data/mimocode/mimo.json" \
    FEISHU_ENABLED=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(str(d.get('feishu_enabled',False)).lower())") \
    FEISHU_APP_ID=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('feishu_app_id',''))") \
    FEISHU_APP_SECRET=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('feishu_app_secret',''))") \
    WECHAT_ENABLED=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(str(d.get('wechat_enabled',False)).lower())") \
    WECHAT_CORP_ID=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('wechat_corp_id',''))") \
    WECHAT_AGENT_ID=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('wechat_agent_id',''))") \
    WECHAT_SECRET=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('wechat_secret',''))") \
    WECHAT_TOKEN=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('wechat_token',''))") \
    WECHAT_ENCODING_AES_KEY=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('wechat_encoding_aes_key',''))") \
    PERSONAL_WECHAT_ENABLED=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(str(d.get('personal_wechat_enabled',False)).lower())") \
    HA_MCP_URL=$(python3 -c "import json; d=json.load(open('/data/options.json')); print(d.get('ha_mcp_url',''))") \
    /usr/bin/python3 /usr/share/mimocode/webui/server.py
