# mimo-code/rootfs/etc/s6-overlay/s6-rc.d/mimocode-webui/run
#!/usr/bin/with-contenv bashio

MIMO_PORT=$(bashio::config 'port')
exec env MIMO_PORT="${MIMO_PORT}" /usr/bin/python3 /usr/share/mimocode/webui/server.py
