# mimo-code/rootfs/etc/s6-overlay/s6-rc.d/mimocode-webui/run
#!/usr/bin/with-contenv bashio

CONFIG_PATH=/data/options.json
if [ -f "$CONFIG_PATH" ]; then
    MIMO_PORT=$(python3 -c "import json; print(json.load(open('$CONFIG_PATH')).get('port', 14096))")
else
    MIMO_PORT=14096
fi

# The web UI talks to the inner (127.0.0.1) port that mimo serve listens on
INNER_PORT=$((MIMO_PORT - 1))

# Persist channel config to /data (survives add-on updates), not the image layer
mkdir -p /data/mimocode

exec env MIMO_PORT="${INNER_PORT}" MIMO_CONFIG="/data/mimocode/mimo.json" \
    /usr/bin/python3 /usr/share/mimocode/webui/server.py
