# mimo-code/rootfs/etc/s6-overlay/s6-rc.d/ha-mcp/finish
#!/usr/bin/with-contenv bashio
# Cleanup on service stop (if needed)
exit 0
